# downloads_agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <!-- Unique identifier for the LaunchAgent -->
    <key>Label</key>
    <string>com.user.downloadsorganizer</string>

    <!-- Command to run -->
    <key>ProgramArguments</key>
    <array>
        <!-- Stable Homebrew Python path -->
        <string>/opt/homebrew/bin/python3</string>

        <!-- Path to the organizer script -->
        <string>/Users/YOURUSERNAME/Scripts/downloads_organizer.py</string>

        <!-- Run once per trigger -->
        <string>--once</string>

        <!-- Enable detailed logging -->
        <string>--verbose</string>

        <!-- Log file location -->
        <string>--logfile</string>
        <string>/Users/YOURUSERNAME/Library/Logs/DownloadsOrganizer.log</string>
    </array>

    <!-- Run once when the user logs in -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Run every 20 seconds -->
    <key>StartInterval</key>
    <integer>20</integer>

    <!-- Stdout / stderr for debugging -->
    <key>StandardOutPath</key>
    <string>/tmp/downloadsorganizer.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/downloadsorganizer.err</string>

</dict>
</plist>